# tlp/99-slstatus-battery.rules
SUBSYSTEM=="power_supply", ACTION=="change", RUN+="/usr/local/bin/slstatus-battery"
